# i3-picture-of-the-day.service
[Unit]
Description=Daily Wallpaper Fetcher

[Service]
ExecStart=/usr/bin/i3-picture-of-the-day.sh
