# i3-picture-of-the-day.timer
[Unit]
Description=Wallpaper Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
